# init/are.timer
[Unit]
Description=run area backup daily

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target