# init/dbns
#!/bin/sh

daemon="/usr/local/bin/dbnsd"
daemon_flags="-p 4000"
rc_reload="NO"
rc_usercheck="NO"

. /etc/rc.d/rc.subr

rc_cmd $1
